# archiso-profile/airootfs/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=/setup.sh
Requires=network-online.target
